# deployment/vibe-webhook.service
[Unit]
Description=Vibe Kanban Webhook Handler
After=network.target

[Service]
Type=simple
User=deployer
WorkingDirectory=/opt/divestreams-automation
ExecStart=/usr/bin/node scripts/vibe-webhook-handler.mjs
Restart=always
RestartSec=10

# Environment variables
Environment=NODE_ENV=production
Environment=PORT=3000
EnvironmentFile=/opt/divestreams-automation/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vibe-webhook

[Install]
WantedBy=multi-user.target
